# systemd/gnome_idle_monitor.service
[Unit]
Description=GNOME Idle Monitor Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 /usr/local/bin/gnome_idle_monitor
Restart=always
RestartSec=5
WorkingDirectory=/tmp

[Install]
WantedBy=default.target
